# docs/com.engram.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.engram.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/ABSOLUTE/PATH/TO/engram/scripts/run_bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/engram</string>

    <!-- Keep the bot running -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if crashed, with throttle -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/engram/logs/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/engram/logs/bot_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
<!--
Installation:
1. Replace all /ABSOLUTE/PATH/TO/engram placeholders with your real project path.
2. Copy this file to ~/Library/LaunchAgents/
3. Load with: launchctl load ~/Library/LaunchAgents/com.engram.bot.plist
4. Check status: launchctl list | grep engram.bot

To stop: launchctl stop com.engram.bot
To unload: launchctl unload ~/Library/LaunchAgents/com.engram.bot.plist
To check logs: tail -f /ABSOLUTE/PATH/TO/engram/logs/bot.log
-->
