# docs/com.engram.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.engram.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/hani/engram/scripts/daily_delta_sync.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/hani/engram</string>

    <!-- Run daily at 6 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/hani/engram/logs/daily_sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hani/engram/logs/daily_sync_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
<!--
Installation:
1. Copy this file to ~/Library/LaunchAgents/
2. Load with: launchctl load ~/Library/LaunchAgents/com.engram.daily.plist
3. Check status: launchctl list | grep engram

To unload: launchctl unload ~/Library/LaunchAgents/com.engram.daily.plist
To run manually: launchctl start com.engram.daily
-->
